# config/inittab
mnt:0:once:/bin/readfstab.lua
shel:1:respawn:/bin/sh.lua
