# config/inittab
mnt:1:wait:/bin/readfstab.lua
shel:1:respawn:/bin/sh.lua
